# docker/supervisord_auth.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:gunicorn]
command=gunicorn --bind 0.0.0.0:8000 daas_py_auth.wsgi:application
autostart=true
autorestart=true
startretries=3                             # Retry 3 times before failing
stderr_logfile=/var/log/gunicorn.err.log
stdout_logfile=/var/log/gunicorn.out.log
redirect_stderr=true                       # Merge stderr into stdout
stdout_logfile_maxbytes=10MB               # Rotate logs when too large
stdout_logfile_backups=5                   # Keep up to 5 backup logs
stopasgroup=true                           # Ensure all child processes stop
killasgroup=true                           # Ensure complete cleanup on stop
